# ansible/roles/gitlab-server/templates/gitlab-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Weekly GitLab backup timer

[Timer]
OnCalendar=Sun *-*-* 03:00:00 UTC
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
