# ansible/roles/config-server/templates/pxe-sync.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Hourly PXE artifact sync from iac-control

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
