# infrastructure/seedbox-vm/qbit-proxy.service
[Unit]
Description=qBittorrent TCP proxy to seedbox VM
After=network.target
[Service]
ExecStart=/usr/bin/socat TCP-LISTEN:18080,fork,reuseaddr TCP:${SEEDBOX_IP}:8080
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
