# ansible/roles/iac-control/templates/gitlab-runner-cache-cleanup.service.j2
# {{ ansible_managed }}
[Unit]
Description=GitLab Runner cache cleanup
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gitlab-runner-cache-cleanup.sh
