# ansible/roles/iac-control/templates/netbox-sync-proxmox.timer.j2
[Unit]
Description=Sync Proxmox VMs to NetBox every {{ netbox_sync_schedule }}

[Timer]
OnCalendar={{ netbox_sync_schedule }}
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
