# ansible/roles/iac-control/templates/proxmox-snapshot.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Daily Proxmox VM Snapshots

[Timer]
OnCalendar={{ proxmox_snapshot_schedule }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
