# ansible/roles/iac-control/templates/unifi-collector.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Run UniFi collector every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
