# ansible/roles/wazuh-server/templates/wazuh-vuln-exporter.service.j2
[Unit]
Description=Wazuh Vulnerability Prometheus Exporter
After=wazuh-indexer.service
Requires=wazuh-indexer.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/sentinel/wazuh-vuln-exporter.py
Restart=on-failure
RestartSec=30
User=root
Group=root

[Install]
WantedBy=multi-user.target
