# ansible/roles/minio-server/templates/b2-backup.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Hourly B2 backup sync timer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
